# systemd/fetch-impressions.timer
[Unit]
Description=Fetch impressions from the collector

[Timer]
OnBootSec=60
OnCalendar=*:0/1

[Install]
WantedBy=facebook-tracking-exposed.target
